# launchd/com.joemezzadri.market-dashboard.fetch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.joemezzadri.market-dashboard.fetch</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/Users/joemezzadri/Documents/market-dashboard/scripts/run-daily-fetch.sh</string>
	</array>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>7</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>
	<key>StandardOutPath</key>
	<string>/Users/joemezzadri/Library/Logs/market-dashboard-fetch.launchd.out.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/joemezzadri/Library/Logs/market-dashboard-fetch.launchd.err.log</string>
	<key>RunAtLoad</key>
	<false/>
</dict>
</plist>
